# systemd/greasewood@.service
[Unit]
Description=greasewood mesh daemon (%i)
Documentation=https://github.com/cschlick/greasewood
After=network-online.target
Wants=network-online.target
# Only run once this membership is configured (create / join writes it).
ConditionPathExists=/etc/greasewood_%i.toml
# Bound the restart loop. Without this, RestartSec=5 never fills systemd's
# default 10s start-limit window (restarts are 5s apart), so a daemon that can't
# start loops FOREVER, invisibly. 5 failures within 2min instead trips the limit
# → the unit enters a visible `failed` state (systemctl status, and gw watch's
# daemon heartbeat goes stale) instead of thrashing silently.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
# gw run creates WireGuard interfaces and edits routing → runs as root.
ExecStart=/usr/local/bin/gw -c /etc/greasewood_%i.toml run
# always, not on-failure: a mesh daemon down is a node silently rotting toward
# credential expiry, whatever stopped it. Seen in the field: a stray
# `sudo killall python` cleanly SIGTERM'd the daemon, on-failure (correctly)
# ignored the clean exit, and the node spent 20h dying politely. Restart=
# never overrides an explicit `systemctl stop`, so deliberate stops still stick.
Restart=always
RestartSec=5
# Liveness watchdog: gw run pings WATCHDOG=1 after every successful reconcile
# (~5s cadence). A daemon that is alive but no longer reconciling — wedged in
# a way no process supervisor can see — misses its pings and is killed +
# restarted after 120s. Pings are a no-op running outside systemd.
WatchdogSec=120
NotifyAccess=main

# --- sandboxing ---------------------------------------------------------
# The daemon runs as root only for CAP_NET_ADMIN (WireGuard + routing). It
# shells out to ip/wg/nft and, when hosts_sync is on, rewrites /etc/hosts.
# These directives keep an RCE in the daemon from owning the host, without
# breaking any of that. Deliberately NOT set:
#   ProtectSystem=strict/full — the daemon writes /etc/hosts (+ its temp and
#     lock siblings in /etc); strict would EROFS them. 'yes' still makes
#     /usr + /boot read-only.
#   ProtectKernelModules — `ip link add type wireguard` may autoload the
#     module on first use; blocking that would break interface creation.
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_ADMIN
ProtectSystem=yes
ProtectHome=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectClock=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
