# systemd/greasewood.path
[Unit]
Description=Watch for greasewood configuration and start the daemon
Documentation=https://gitlab.com/cschlick/greasewood

[Path]
# Start greasewood.service once /etc/greasewood.toml exists. After a
# config-changing re-join: systemctl restart greasewood.
PathExists=/etc/greasewood.toml
Unit=greasewood.service

[Install]
WantedBy=paths.target
